# docker-dist/etc/supervisor.d/cron.ini
[program:crond]
command=crond -f -s -P
directory=/www
user=root
redirect_stderr=true
autostart=%(ENV_SPV_CRON)s
autorestart=%(ENV_SPV_CRON)s

