# docker-dist/etc/supervisor.d/entrypoint.ini
[program:entrypoint]
command=/www/app/docker-entrypoint.sh
directory=/www
user=root
autorestart=false
exitcodes=0
redirect_stderr=true
startsecs=0
autostart=%(ENV_SPV_ENTRYPOINT)s
