# docker-dist/etc/supervisor.d/nginx.ini
[program:nginx]
command=nginx -g 'daemon off;'
autostart=%(ENV_SPV_NGINX)s
autorestart=%(ENV_SPV_NGINX)s

