# docker-dist/etc/supervisor.d/websocket.ini
[program:websocket]
command=/bin/websocketd --port=7777 --dir=/www/app/websocket
directory=/www
user=root
redirect_stderr=true
autostart=%(ENV_SPV_WEBSOCKET)s
autorestart=%(ENV_SPV_WEBSOCKET)s
